# ops/launchd/local.prism.laptop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  PRISM / lab-scheduler — LaunchAgent for the LAPTOP deploy
  (user vishvajeetn, working copy at ~/Documents/Scheduler/Main).

  Hosts PRISM as a persistent local service on the dev laptop
  itself, independent of the Mac mini + Tailscale path. Listens
  only on 127.0.0.1:5055 by default (see .env) so nothing leaves
  the machine.

  Sibling of ops/launchd/local.prism.plist, which targets the mini.
  Same Label (`local.prism`) — only one of the two plists should
  ever be bootstrapped in a given gui/<uid> domain. The installer
  at scripts/install_launchd.sh picks the correct plist by probing
  $PWD, so a fresh clone onto either machine auto-routes.

  Install:
    ./scripts/install_launchd.sh

  Uninstall:
    launchctl bootout gui/$(id -u)/local.prism
    rm ~/Library/LaunchAgents/local.prism.plist

  Tail logs:
    tail -f ~/Documents/Scheduler/Main/logs/server.log

  Notes:
    * KeepAlive=true restarts the process on any exit. ThrottleInterval
      caps crash loops at one restart per 10s.
    * RunAtLoad=true starts the service at user login, so PRISM is
      already up the moment you unlock the laptop.
    * .env is sourced by start.sh (see that script) so SECRET_KEY +
      DEMO_MODE + CSRF flags survive restarts without leaking into
      launchd's own EnvironmentVariables dict.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.prism</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/vishvajeetn/Documents/Scheduler/Main/scripts/start.sh</string>
        <string>--service</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/vishvajeetn/Documents/Scheduler/Main</string>

    <key>StandardOutPath</key>
    <string>/Users/vishvajeetn/Documents/Scheduler/Main/logs/server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/vishvajeetn/Documents/Scheduler/Main/logs/server.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Interactive</string>

    <!-- Force the Werkzeug reloader OFF under launchd. The reloader
         forks a grandchild and exits the parent, which makes launchd
         mark the service EX_CONFIG-crashed. LAB_SCHEDULER_AUTORELOAD=0
         is the explicit opt-out; see app.py's auto_reload block. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>LAB_SCHEDULER_AUTORELOAD</key>
        <string>0</string>
    </dict>
</dict>
</plist>
