# chooser/launchd/local.catalyst.ravikiran.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Ravikiran ERP — household / estate live lane on the iMac.
  Runs on port 5057 and must boot only from the isolated
  /ERP-Instances/ravikiran-erp/live runtime root.

  Install:
    cp local.catalyst.ravikiran.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) \
      ~/Library/LaunchAgents/local.catalyst.ravikiran.plist
    launchctl kickstart -k gui/$(id -u)/local.catalyst.ravikiran

  Tail logs:
    tail -f ~/ERP-Instances/ravikiran-erp/live/data/logs/server.log
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.catalyst.ravikiran</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/nv/ERP-Instances/ravikiran-erp/live/run.sh</string>
    <string>--service</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/nv/ERP-Instances/ravikiran-erp/live</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/nv/ERP-Instances/ravikiran-erp/live/data/logs/server.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nv/ERP-Instances/ravikiran-erp/live/data/logs/server.log</string>

  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>60</integer>
  <key>ProcessType</key><string>Interactive</string>
</dict>
</plist>
